# quant_system/stock_monitor.service
[Unit]
Description=Stock Portfolio Real-time Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/autostock/quant_system
ExecStart=/usr/bin/python3 /root/autostock/quant_system/realtime_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
